# App/crontab/crontab
# START CRON JOB
*/10 * * * * python3 /app/MoveVideo.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB