# App/crontab/crontab
*/30 * * * * python /app/MoveVideo.py >> /logfile